# config/opencode-web.service
[Unit]
Description=OpenCode Web Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=/home/linuxbrew/.linuxbrew/bin/opencode web --port 4096 --hostname 0.0.0.0
Restart=on-failure
RestartSec=5
Environment=OPENCODE_SERVER_USERNAME=matvey
Environment=OPENCODE_SERVER_PASSWORD={{OPENCODE_SERVER_PASSWORD}}

[Install]
WantedBy=default.target
